# services/airflow-flower.service
[Unit]
Description=airflow-flower

[Service]
Type=simple
User=root
Restart=on-failure
LimitNOFILE=102400

ExecStart=/data/software/app/airflow/start-flower.sh

[Install]
WantedBy=multi-user.target
